# debian/postinst
#!/bin/sh

set -e

#DEBHELPER#

if [ "$1" = configure ]; then
    pov-update-server-page -q && invoke-rc.d apache2 reload || true
fi
